# etc/systemd/system/docker.service.d/storage.conf
[Service]
Environment="DOCKER_DATA_DIR=/home/docker"
ExecStart= 
ExecStart=/usr/bin/dockerd --data-root=${DOCKER_DATA_DIR} -H fd://

