# etc/systemd/system/docker.service.d/storage.conf
[Service]
ExecStart= 
ExecStart=/usr/bin/dockerd --data-root={{DOCKER_DATA_ROOT}} -H fd://

