# system/etc/systemd/system/docker.service.d/storage.conf
[Service]
ExecStart= 
ExecStart=/usr/bin/dockerd --data-root=/media/b/var/docker/ -H fd://

